# firecracker@.service
# firecracker@.service — one instance per VM (firecracker@<name>)
# Enabled instances are restored after a host reboot, and the supervisor owns
# the VM lifecycle: guest reboots and kernel panics auto-recover (Restart=always).
[Unit]
Description=Firecracker microVM %i
After=vps-network.service
Requires=vps-network.service
# Boot-loop backstop: a VM that panics every boot trips the burst limit after
# ~10 rapid restarts and enters `failed` instead of looping forever.
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
ExecStartPre=/usr/bin/clawvps net-up %i
ExecStart=/usr/sbin/firecracker --api-sock /run/vps/%i.sock --config-file /var/lib/vms/%i/config.json
# Attempt a graceful shutdown (x86_64 only); falls through to SIGTERM on failure.
ExecStop=/usr/bin/clawvps send-shutdown %i
ExecStopPost=/usr/bin/clawvps net-down %i
TimeoutStopSec=15
# Always restart: firecracker exits 0 for ALL guest events (poweroff/reboot/panic),
# so on-failure would never fire for a crash. always recovers reboots and panics;
# in-guest poweroff also restarts (indistinguishable at the exit-code level).
# The intentional-down path is `clawvps stop` (a systemd-initiated stop, unaffected).
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
